# scripts/com.openclaw.skill-install.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.skill-install</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/chenggl/workspace/scripts/skill-install.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/chenggl/workspace</string>
    
    <key>StandardOutPath</key>
    <string>/Users/chenggl/workspace/logs/launchd-skill-install.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/chenggl/workspace/logs/launchd-skill-install-error.log</string>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>5</integer>
    </dict>
    
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
